# install/mt-billing-auto-update.timer
[Unit]
Description=Poll GitHub for MT-Billing updates every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
